# examples/systemd/activemq_exporter.socket
[Unit]
Description=ActiveMQ Exporter

[Socket]
ListenStream=9168

[Install]
WantedBy=sockets.target
